# ansible/launchd/com.frias.orchestrator-fleet-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.frias.orchestrator-fleet-monitor</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/python3</string>
    <string>/Users/friasc/.pi/agent/git/github.com/carlosfrias/pi-a2a-communication/scripts/orchestrator-fleet-monitor.py</string>
  </array>
  <key>StartInterval</key>
  <integer>120</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/friasc/Library/Logs/orchestrator-fleet-monitor.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/friasc/Library/Logs/orchestrator-fleet-monitor.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ANSIBLE_VAULT_PASSWORD_FILE</key>
    <string></string>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
