# ansible/systemd/fleet-a2a-health-reporter.service
[Unit]
Description=Fleet A2A Health Reporter
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=friasc
Group=friasc
ExecStart=/home/friasc/fleet-scripts/fleet-health-reporter.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fleet-a2a-health-reporter

# Ensure nvm is available for pi if the reporter invokes it
Environment=PATH=/home/friasc/.nvm/versions/node/v24.15.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=NVM_DIR=/home/friasc/.nvm
