# ansible/systemd/ollama-idle-unload.service
[Unit]
Description=Ollama Idle Model Unload (calls fleet-health-check --ollama-only)
After=ollama.service

[Service]
Type=oneshot
User=friasc
Group=friasc
ExecStart=/home/friasc/fleet-scripts/fleet-health-check.sh --ollama-only
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama-idle-unload

# Ensure nvm is available for pi
Environment=PATH=/home/friasc/.nvm/versions/node/v24.15.0/bin:/usr/local/bin:/usr/bin:/bin
Environment=NVM_DIR=/home/friasc/.nvm